# etc-systemd-system/backupmysystemusb2usb.service
[Unit]
After=ssh.service

[Service]
ExecStart=/home/stephane/Code/Python/backupmysystemusb2usb/run.sh

[Install]
WantedBy=default.target
