# etc-systemd-system/backupmysystemusb2usb.service
[Unit]
After=ssh.service

[Service]
WorkingDirectory=/home/stephane/Code/Python/backupmysystemusb2usb/
ExecStart=/home/stephane/Code/Python/backupmysystemusb2usb/run.py config.yml

[Install]
WantedBy=default.target
